# description/root/etc/systemd/system/sshd.service.d/10-lastlog.conf
[Service]
ExecStartPre=/usr/bin/touch /var/log/btmp
ExecStartPre=/usr/bin/touch /var/log/lastlog
